# tpm-fido-kde/99-tpm-custom.rules
KERNEL=="tpmrm[0-9]*", RUN+="/bin/setfacl -m g:_tpmfido:rw $env{DEVNAME}"
