# tpm-fido-kde/99-uhid.rules
KERNEL=="uhid", SUBSYSTEM=="misc", RUN+="/bin/setfacl -m g:_tpmfido:rw $env{DEVNAME}"
