# automount-lvm-snapshots/mnt-snapshots-www.automount
[Unit]
Description=www snapshots automount

[Automount]
Where=/mnt/snapshots/www
TimeoutIdleSec=60s
